# proxy/supervisord.conf
#[unix_http_server]
#file=/logs/supervisor.sock
#
#[rpcinterface:supervisor]
#supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
#
#[supervisorctl]
#serverurl=unix:///logs/supervisor.sock

[supervisord]
pidfile = /tmp/supervisord.pid
logfile = /logs/supervisord.log

[program:zabbix_server]
autorestart=true
command=/usr/sbin/zabbix_proxy -f
